# mycardtracker-api.service
[Unit]
Description=MyCardTracker API Server
After=network.target

[Service]
User=webuser
WorkingDirectory=/var/www/mycardtracker/dist/api
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=mycardtracker-api
Environment=NODE_ENV=production
Environment=PORT=8080

[Install]
WantedBy=multi-user.target 